# deploy/systemd/porteight-camera-ai.service
[Unit]
Description=PortEight Camera AI (RTMP ingest + HLS + recordings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=porteight
Group=porteight
WorkingDirectory=/opt/porteight-camera-ai

# Persist DB + recordings + HLS under /var/lib (change as needed)
Environment=PERSIST_DIR=/var/lib/porteight-camera-ai
Environment=WEB_DIR=/opt/porteight-camera-ai/web
Environment=GIN_MODE=release

# Shared auth with Suvidhi
# - NEXTAUTH_SECRET must match the secret used by Suvidhi's NextAuth/Auth.js deployment.
# - AUTH_COOKIE_SECURE should be "false" for plain HTTP on :8080 and "true" behind HTTPS.
# - SIGNING_SECRET is used for signed playback/download URLs.
Environment=NEXTAUTH_SECRET=CHANGE_ME_TO_MATCH_SUVIDHI_NEXTAUTH_SECRET
Environment=AUTH_COOKIE_SECURE=false
Environment=SIGNING_SECRET=CHANGE_ME_TO_A_LONG_RANDOM_SIGNING_SECRET

# If you want to override individually:
# Environment=DB_PATH=/var/lib/porteight-camera-ai/camera_ai.db
# Environment=RECORDINGS_DIR=/var/lib/porteight-camera-ai/recordings
# Environment=HLS_DIR=/var/lib/porteight-camera-ai/hls

ExecStart=/opt/porteight-camera-ai/server
Restart=always
RestartSec=2

# Hardening (safe defaults; loosen if needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/porteight-camera-ai

[Install]
WantedBy=multi-user.target
