# deploy/systemd/porteight-camera-ai.service
[Unit]
Description=PortEight Camera AI (RTMP ingest + HLS + recordings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=porteight
Group=porteight
WorkingDirectory=/opt/porteight-camera-ai

# Persist DB + recordings + HLS under /var/lib (change as needed)
Environment=PERSIST_DIR=/var/lib/porteight-camera-ai
Environment=WEB_DIR=/opt/porteight-camera-ai/web
Environment=GIN_MODE=release

# Session cookie config
# - SESSION_AUTH_KEY must be 32+ bytes (required)
# - SESSION_SECURE should be "false" if you are serving plain HTTP on :8080, and "true" if behind HTTPS.
Environment=SESSION_AUTH_KEY=CHANGE_ME_TO_A_LONG_RANDOM_SECRET_32+_BYTES
Environment=SESSION_SECURE=false

# If you want to override individually:
# Environment=DB_PATH=/var/lib/porteight-camera-ai/camera_ai.db
# Environment=RECORDINGS_DIR=/var/lib/porteight-camera-ai/recordings
# Environment=HLS_DIR=/var/lib/porteight-camera-ai/hls

ExecStart=/opt/porteight-camera-ai/server
Restart=always
RestartSec=2

# Hardening (safe defaults; loosen if needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/porteight-camera-ai

[Install]
WantedBy=multi-user.target
